# luci-app-easy-frpc/root/etc/uci-defaults/80_easy-frpc
#!/bin/sh

touch /etc/config/easy-frpc
uci set easy-frpc.first=first
uci set easy-frpc.second=second
uci set easy-frpc.third=third
uci commit

return 0